# install/blueeye-agent.service
[Unit]
Description=BlueEye Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=blueeye
WorkingDirectory=/opt/blueeye-agent
EnvironmentFile=/etc/blueeye-agent/env
ExecStart=/usr/bin/node /opt/blueeye-agent/index.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
