# install/blueeye-agent.service
[Unit]
Description=BlueEye Network Agent
Documentation=https://github.com/gudlever-lgtm/blueeye-agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Adjust paths for the target device (e.g. a Cisco switch guest shell).
WorkingDirectory=/opt/blueeye-agent
EnvironmentFile=/opt/blueeye-agent/.env
ExecStart=/usr/bin/python3 /opt/blueeye-agent/main.py
Restart=always
RestartSec=10
User=blueeye

[Install]
WantedBy=multi-user.target
